# scripts/deploy/contract-review.service
# scripts/deploy/contract-review.service
#
# Systemd unit file for the contract review FastAPI server.
# Copy to /etc/systemd/system/contract-review.service
#
# Commands:
#   sudo systemctl daemon-reload
#   sudo systemctl start contract-review
#   sudo systemctl enable contract-review   # auto-start on boot
#   sudo systemctl status contract-review
#   journalctl -u contract-review -f        # live logs

[Unit]
Description=Contract Review FastAPI (uvicorn)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/contract-review

# All API keys are loaded from this file — never hardcode them above.
EnvironmentFile=/home/ubuntu/contract-review/.env

# Single worker is REQUIRED:
#   - Each worker loads the sentence-transformers model independently (~4 min
#     first-time cost per process).
#   - In-memory job state (_jobs dict, asyncio.Queues used for SSE) is not
#     shared across workers, so an SSE stream connection may land on a
#     different worker than the one running the job.
# Scale up by running multiple machines behind a load balancer with
# sticky-session routing instead.
ExecStart=/home/ubuntu/contract-review/venv/bin/uvicorn \
    api.routes:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1

Restart=always
RestartSec=5
# Startup blocks on a synchronous embedder prewarm (3-4 min on cold caches).
# Default 90s timeout would kill us mid-load.
TimeoutStartSec=600
TimeoutStopSec=30

# Log to journald — read with: journalctl -u contract-review -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
